# roles/nut/templates/etc_udev_rules.d_99-nut-ups.rules.j2
# {{ ansible_managed }}
# Unbind kernel usbfs/usbhid driver from UPS HID device so NUT's usbhid-ups can claim it.
# Without this, NUT gets "Can't claim USB device" because the kernel driver holds the interface.

# ACTION=="add", SUBSYSTEM=="usb", ATTR{idVendor}=="{{ nut_ups_vendorid }}", ATTR{idProduct}=="{{ nut_ups_productid }}", MODE="0660", GROUP="nut", RUN+="/bin/sh -c 'for iface in /sys$devpath:*; do [ -e \"$iface/driver/unbind\" ] && echo $(basename $iface) > $iface/driver/unbind; done; true'"

# ─── Stable symlinks for two identical UPS units (commented out) ───────
# Two identical UPS units share the same vendor:product ID. Without stable
# symlinks, /dev/bus/usb paths can swap after reboot.
# udev rules with KERNELS (= physical USB port path) solve this.
#
# Determine KERNELS path:
#   udevadm info -a /dev/bus/usb/001/003 | grep KERNELS
#   → e.g. KERNELS=="1-2" (USB port) or KERNELS=="1-2.1" (via hub)
#
# Uncomment and adjust KERNELS values:
#

# elite800: port 1-4 = greencell (600W), port 1-1 = greencell2 (2000VA 1200W)
ACTION=="add", SUBSYSTEM=="usb", ATTR{idVendor}=="{{ nut_ups_vendorid }}", ATTR{idProduct}=="{{ nut_ups_productid }}", KERNELS=="1-4", SYMLINK+="ups-1", MODE="0660", GROUP="nut", RUN+="/bin/sh -c 'for iface in /sys$devpath:*; do [ -e \"$iface/driver/unbind\" ] && echo $(basename $iface) > $iface/driver/unbind; done; true'"
ACTION=="add", SUBSYSTEM=="usb", ATTR{idVendor}=="{{ nut_ups_vendorid }}", ATTR{idProduct}=="{{ nut_ups_productid }}", KERNELS=="1-1", SYMLINK+="ups-2", MODE="0660", GROUP="nut", RUN+="/bin/sh -c 'for iface in /sys$devpath:*; do [ -e \"$iface/driver/unbind\" ] && echo $(basename $iface) > $iface/driver/unbind; done; true'"

#
# Then switch to symlink-based ports in defaults/main.yml:
#   nut_ups_port: "/dev/ups-1"
#   nut_ups_secondary_port: "/dev/ups-2"
# and replace the generic rule above with the two specific ones.
# ────────────────────────────────────────────────────────────────────────
