# roles/dgx_prepare/templates/qsfp_sriov.service.j2
# Managed by Ansible — do not edit
[Unit]
Description=Enable {{ qsfp_sriov_numvfs }} SR-IOV VFs on QSFP PF {{ qsfp_sriov_pf_pci }}
DefaultDependencies=no
Before=network-pre.target
After=systemd-modules-load.service

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/bin/bash -c 'echo {{ qsfp_sriov_numvfs }} > /sys/bus/pci/devices/{{ qsfp_sriov_pf_pci }}/sriov_numvfs'
{% for vf_idx in range(0, qsfp_sriov_numvfs | int) %}
ExecStartPost=/bin/bash -c 'ip link set {{ qsfp_sriov_vf_prefix }}v{{ vf_idx }} down && ip link set {{ qsfp_sriov_vf_prefix }}v{{ vf_idx }} address {{ '02:00:0a:0a:%02x:%02x' | format(host_octet | int, vf_idx) }} && ip link set {{ qsfp_sriov_vf_prefix }}v{{ vf_idx }} up'
{% endfor %}
ExecStop=/bin/bash -c 'echo 0 > /sys/bus/pci/devices/{{ qsfp_sriov_pf_pci }}/sriov_numvfs'

[Install]
WantedBy=multi-user.target
